# gateway-autofix/ai.openclaw.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.watchdog</string>

    <key>Comment</key>
    <string>OpenClaw Gateway Watchdog — monitors health, triggers Claude Code auto-fix on repeated crashes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/clekee/work/openclaw-skills/gateway-autofix/scripts/openclaw-watchdog.sh</string>
    </array>

    <!-- Run every 60 seconds -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/clekee/.openclaw/logs/watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/clekee/.openclaw/logs/watchdog.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/clekee</string>
        <key>PATH</key>
        <string>/Users/clekee/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>OPENCLAW_GATEWAY_PORT</key>
        <string>18789</string>
        <!-- Uncomment to enable Telegram notifications:
        <key>OPENCLAW_FIX_TELEGRAM_TARGET</key>
        <string>957022683</string>
        -->
    </dict>
</dict>
</plist>
